# Arch_PKGBUILDs/postmarketos-firmware-loader/50-firmware.rules
# do not edit this file, it will be overwritten on update
SUBSYSTEM=="firmware", ACTION=="add", RUN+="/usr/lib/firmwareload.sh"
