# Arch_PKGBUILDs/postmarketos-usb-networking/usb_networking.service
[Unit]
After=sysinit.target

[Service]
Type=simple
ExecStart=/usr/lib/danctnix/start_usb_networking

[Install]
WantedBy=multi-user.target
