# Arch_PKGBUILDs/postmarketos-wcnss-wlan/wcnss-wlan.service
[Unit]
After=sysinit.target

[Service]
Type=simple
ExecStart=/usr/lib/danctnix/start_wcnss_wlan

[Install]
WantedBy=multi-user.target
